# backup-service/backup.timer
[Unit]
Description=Backup dotfiles and packages

[Timer]
OnBootSec=60min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
